# ros2-humble/supervisord.conf
[supervisord]
nodaemon=true

[program:xvfb]
command=Xvfb :99 -screen 0 1920x1080x24
autorestart=true

[program:x11vnc]
command=x11vnc -display :99 -nopw -forever -shared
autorestart=true